# dot_Library/LaunchAgents/com.mickey-kras.bw-session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mickey-kras.bw-session</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>BW=$("$HOME/.local/bin/bw-gate-helper" get 2>/dev/null); if [ -n "$BW" ]; then /bin/launchctl setenv BW_SESSION "$BW"; fi</string>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/bw-session-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bw-session-agent.log</string>
</dict>
</plist>
